# templates/wal-g-prometheus-exporter.service.epp
[Unit]
Description=Wal-g prometheus exporter

[Service]
Restart=always
User=postgres
Group=postgres
EnvironmentFile=/usr/local/bin/exporter.env
ExecStart=/usr/local/bin/wal-g-prometheus-exporter --debug <%= $datadir %>/pg_wal/archive_status/
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
